# zenplay.service
[Unit]
Description=zenplay
After=redis.service

[Service]
Type=simple
ExecStart=/usr/local/bin/zenplay

[Install]
WantedBy=multi-user.target
